# 6-systemd/monitor.service
[Unit]
Description=Monitor

[Service]
Type=simple
EnvironmentFile=-/etc/sysconfig/monitor
ExecStart=/usr/local/bin/monitor.sh
